# root/lib/systemd/system/horst_deployed.timer
[Unit]
Description=HORST DEPLOYED Signal Timer

[Timer]
OnBootSec=45min
Unit=horst_deployed.service

[Install]
WantedBy=LEOP.target
