# root/etc/systemd/system/horst_detumbler.timer
[Unit]
Description=horst detumble timer
[Timer]
OnBootSec=6min
Unit=horst_detumbler.service
[Install]
WantedBy=LEOP.target
