# root/etc/systemd/system/horst_detumbler.timer
[Unit]
Description=HORST DEPLOYED signal timer

[Timer]
OnBootSec=45min
Unit=horst_detumbler.service

[Install]
WantedBy=LEOP.target
