# root/etc/systemd/system/horst_deployed.timer
[Unit]
Description=HORST DEPLOYED signal timer

[Timer]
OnBootSec=45min
Unit=horst_deployed.service

[Install]
WantedBy=LEOP.target
